# deploy/forms.conf
[program:forms]
command = /home/ubuntu/deploy/deploy
user = root
