# deploy/forms.conf
[program:forms]
command = /home/ubuntu/deploy/deploy.sh
user = root
